# service/sapping.service
[Unit]
Description=Hiding sapservices file from sapinit
Documentation=man:sapping(8)
#Documentation=man:sapservices-move(8)
Before=sapinit.service pacemaker.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/sapping
#ExecStart=/usr/sbin/sapservices-move --away
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target pacemaker.service
